# presence-ssh-mqtt/rootfs/etc/services.d/unifi-ssh-presence/finish
#!/usr/bin/with-contenv bash
echo "[unifi-ssh-presence] stopped. exit_code=$1 signal=$2"
exit 0
