# ha-moode-radios-addon/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
set -euo pipefail

export PORT="${PORT:-7860}"

exec /opt/venv/bin/python /app/main.py
