# stremio-hell-ha/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
set -euo pipefail

exec php-fpm82 -F
