# stremio-sos-ha/rootfs/etc/services.d/sosac-stremio/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "Sosac Stremio service stopped."
