# stremio-prehrajto-addon/rootfs/etc/services.d/stremio-prehrajto/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

export PORT="7654"
export BASE_URL="$(bashio::config 'base_url')"
export TMDB_API_KEY="$(bashio::config 'tmdb_api_key')"
export TMDB_LANGUAGE="$(bashio::config 'tmdb_language')"
export PREHRAJTO_BASE="$(bashio::config 'prehrajto_base')"
export PREHRAJTO_USER_AGENT="$(bashio::config 'prehrajto_user_agent')"
export DEFAULT_SEARCH_LIMIT="$(bashio::config 'default_search_limit')"
export DEFAULT_PREMIUM="$(bashio::config 'default_premium')"
export CACHE_TTL_SECONDS="$(bashio::config 'cache_ttl_seconds')"
export LOG_LEVEL="$(bashio::config 'log_level')"
export CONFIG_SECRET="$(bashio::config 'config_secret')"
export DATA_DIR="/data"

bashio::log.info "Starting Stremio Prehraj.to add-on on port ${PORT}"
exec node /app/src/index.js
