# stremio-hell-ha/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bash
set -euo pipefail

PORT="${PORT:-80}"

cd /var/www/html
exec php -S 0.0.0.0:${PORT} index.php
